# rancher/cron/crontab
# for zabbix user
*/10 *  * * * /home/sup/code/bash_configs/rancher/cron/git-pull.sh

